# systemd/logibar-hidpp-monitor.service
[Unit]
Description=Logibar HID++ Battery Monitor (Keyboard/Mouse)
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/logibar-hidpp-monitor
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
