# systemd/logibar-headset-monitor.service
[Unit]
Description=Logibar PRO X 2 Headset Battery Monitor
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/logibar-headset-monitor
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
